# templates/puppet-data-connector-enhancer.service.epp
<%- | String $script_path,
      String $dropzone_file | -%>
[Unit]
Description=Puppet Data Connector Enhancer
After=network.target

[Service]
Type=oneshot
User=pe-puppet
Group=pe-puppet
ExecStart=<%= $script_path %> -q -o <%= $dropzone_file %>
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target